# metarmap.service
[Unit]
Description=LED Metar Map
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/pi/METARMap/.env/bin/python3 /home/pi/METARMap/metar.py

[Install]
WantedBy=multi-user.target
